# build/mfa/runit/app/run
#!/bin/sh
exec chpst -u app node /app/bin/src/server.js --config /app/conf/mfa.yml
