# build/service-mfa/runit/app/run
#!/bin/sh
exec chpst -u app node /app/bin/server.js --config /app/conf/service-mfa.yml
